# systemd/system/update_background.timer
[Unit]
Description=Run update_background every minute

[Timer]
OnBootSec=1min
OnActiveSec=1min

[Install]
WantedBy=timers.target
