# systemd/system/update_background.service
Unit]
Description=Run update_background

[Service]
ExecStart=/bin/bash /home/attila/bin/update_background
User=attila
